# alpine/alpine-root/etc/periodic/hourly/reload-cloudflare-ips.sh
#!/usr/bin/env sh

#update-cloudflare-ips.sh && logger 'received new Cloudflare IPs' && /etc/init.d/local restart | logger
