# alpine/alpine-root/etc/local.d/00-ipv6-tunnels.stop
#!/bin/sh

busybox ip -6 route flush dev || :
for i in netassist ip4market he ; do
  busybox ip link set dev "$i" down || :
  busybox ip tunnel del "$i" || :
done
