# alpine/alpine-root/etc/local.d/03-ssh-antibruteforce.stop
#!/usr/bin/env bash

kill $(ps uax | grep -E '(03-ssh-antibruteforce.start|tail -n1 -F /var/log/messages)' | grep -v grep | awk '{print $2}')
