# alpine/alpine-root/etc/local.d/03-syslog-handlers.stop
#!/usr/bin/env bash

kill $(ps uax | grep "tail -n1 -F ${SYSLOG}" | grep -v grep | awk '{print $2}')
