# alpine/alpine-root/etc/local.d/03-syslog-handlers.stop
#!/usr/bin/env bash

kill $(ps uax | grep 'tail -n1 -F /tmp/messages' | grep -v grep | awk '{print $2}')
