# alpine/alpine-root/etc/local.d/local.start
#!/bin/sh

chvt 12

setterm -blank 0
setterm -powersave off
setterm -powerdown 0

update-allowlist-ips.sh || :
