# alpine/alpine-root/etc/init.d/quicssh-rs
#!/sbin/openrc-run

name="quicssh-rs"

command="/usr/local/bin/${name}"
command_user="nobody:nobody"
command_background="true"
command_args="${quicssh_rs_args}"

capabilities="^cap_net_bind_service"
pidfile="/run/${name}.pid"
output_logger="logger"
error_logger="logger"

depend() {
	need net
	use logger
}
