# alpine/alpine-root/etc/init.d/quicssh-rs
#!/sbin/openrc-run

command="/usr/local/bin/${RC_SVCNAME}"
command_user="nobody:nobody"
command_background="true"
command_args="${quicssh_rs_args}"

capabilities="^cap_net_bind_service"
pidfile="/run/${RC_SVCNAME}.pid"
output_logger="logger -t ${RC_SVCNAME}"
error_logger="logger -t ${RC_SVCNAME}"

depend() {
	need net
	use logger
}
